# debian/postinst
#!/bin/bash -e

systemctl disable lxc
systemctl disable lxc-net

exit 0